# defenra-agent-healthcheck.service
[Unit]
Description=Defenra Agent Health Check
After=defenra-agent.service
Requires=defenra-agent.service

[Service]
Type=oneshot
ExecStart=/opt/defenra-agent/defenra-agent-healthcheck.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=defenra-agent-healthcheck
